# com.example.pm-news-pulse-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.yourname.pm-news-pulse-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/path/to/pm-news-pulse/sync_to_vault.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/path/to/pm-news-pulse</string>
    <key>RunAtLoad</key>
    <true/>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/path/to/pm-news-pulse/logs/sync.log</string>
    <key>StandardErrorPath</key>
    <string>/path/to/pm-news-pulse/logs/sync-error.log</string>
</dict>
</plist>
